# docker/cron.d/cleanup
* * * * * root find /opt/calibre-api/files -type f -mtime +1 -delete
