# docker/supervisord.conf
[supervisord]
nodaemon = true
logfile = /dev/null
logfile_maxbytes = 0
pidfile = /run/supervisord.pid
user=root

[program:cron]
command = /usr/sbin/cron -fL 15
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes=0
user = root
autostart = true
autorestart = true
priority = 20

[program:calibre-api]
directory=/opt/calibre-api
command = /opt/calibre-api/calibre-api
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes=0
user = nobody
autostart = true
autorestart = true
priority = 5
